# etc/service/mediatomb/run
#!/bin/sh

. /etc/container_environment.sh

exec mediatomb -p 8000 -a /srv/mediatomb
